# misc/atune-engine.service
[Unit]
Description=A-Tune AI service
After=network.target
Requires=polkit.service

[Service]
ExecStart=python3 /usr/libexec/atuned/analysis/app_engine.py /etc/atuned/engine.cnf

[Install]
WantedBy=multi-user.target
